# templates/prometheus.conf.j2
# Run prometheus

start on startup

chdir {{prometheus_server_dir}}

script
   ./prometheus -config.file /etc/prometheus/prometheus.yml
end script
